# scripts/galaxea/r1lite/host/dimos-r1lite.service
[Unit]
Description=DimOS R1 Lite current release
Requires=docker.service dimos-hostnet.service
After=docker.service dimos-hostnet.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=exec
WorkingDirectory=/opt/dimos/current
ExecStart=/usr/bin/docker compose up --no-color --abort-on-container-exit
ExecStopPost=/usr/bin/docker compose down --timeout 35
Restart=on-failure
RestartSec=10
TimeoutStopSec=45

[Install]
WantedBy=multi-user.target
